# launchd/com.propertyagent.monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <!-- Identity -->
  <key>Label</key>
  <string>com.propertyagent.monthly</string>

  <!-- Script to run — shell wrapper activates the venv before calling Python -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/william/Desktop/ClaudeCode/run_monthly.sh</string>
  </array>

  <!-- Working directory (launchd does not inherit a shell cwd) -->
  <key>WorkingDirectory</key>
  <string>/Users/william/Desktop/ClaudeCode</string>

  <!-- Schedule: 1st of every month at 07:00 local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Do NOT run immediately when the plist is loaded -->
  <key>RunAtLoad</key>
  <false/>

  <!-- Log files in the project logs/ directory -->
  <key>StandardOutPath</key>
  <string>/Users/william/Desktop/Property Research/logs/monthly.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/william/Desktop/Property Research/logs/monthly.stderr.log</string>

  <!-- Restart on crash with a 60-second throttle -->
  <key>ThrottleInterval</key>
  <integer>60</integer>

  <!-- Inherit a usable PATH that includes /usr/local/bin and Homebrew -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/william</string>
  </dict>

</dict>
</plist>
